# web/deployment/trackwise.service
# Raspberry Pi systemd service file for TrackWise Web
#
# Installation:
#   sudo cp trackwise.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable trackwise
#   sudo systemctl start trackwise
#   sudo systemctl status trackwise
#
# Logs:
#   sudo journalctl -u trackwise -f

[Unit]
Description=TrackWise Web — GPX Waypoint Finder
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/trackwise/web/backend
ExecStart=/home/pi/trackwise/web/.venv/bin/uvicorn app:app --host 127.0.0.1 --port 8000 --workers 1
Restart=always
RestartSec=5

# Environment
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONDONTWRITEBYTECODE=1

# Resource limits (conservative for Pi)
MemoryMax=400M
CPUQuota=80%

[Install]
WantedBy=multi-user.target
